# systemd/dae-api.service
[Unit]
Description=DAE API

[Service]
WorkingDirectory=/root
ExecStart=/usr/bin/node /root/dae-api/server.js

[Install]
WantedBy=multi-user.target
