# systemd/mongod.service
[Unit]
Description=MongoDB

[Service]
ExecStart=/usr/bin/mongod --dbpath /data

[Install]
WantedBy=multi-user.target
